# debian/dehydrated.cron.d
0 9 * * 6 root /usr/bin/dehydrated -c -f /etc/dehydrated/config
